# docker/app/supervisord.conf
; ERP 全合一容器进程管理:daphne(ASGI) + celery worker + beat + 升级进度中转 + nginx。
; 首启迁移/bootstrap 由 entrypoint 在 supervisord 之前完成。
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:daphne]
command=daphne -b 127.0.0.1 -p 8000 config.asgi:application
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery-worker]
command=celery -A config worker -l info --concurrency=2
directory=/app
autostart=true
autorestart=true
startsecs=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery-beat]
command=celery -A config beat -l info
directory=/app
autostart=true
autorestart=true
startsecs=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:upgrade-relay]
command=python manage.py upgrade_progress_relay
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
